# template/extra/distributionscripts/arch/__SHORTNAME__.service
#copy this to /usr/lib/systemd/system/__SHORTNAME__.service
#systemctl enable __SHORTNAME__
#systemctl start __SHORTNAME__

[Unit]
Description=__SHORTDESCRIPTION__

[Service]
Type=forking
PIDFile=/run/__EXECUTABLE__.pid
ExecStart=/usr/bin/__EXECUTABLE__ --pidfile=/run/__EXECUTABLE__.pid

[Install]
WantedBy=multi-user.target
